# packages/frp/systemd/frpc.service
[Unit]
Description=frpc client
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/frpc -c /etc/frp/frpc.ini

[Install]
WantedBy=multi-user.target

